# minecraft.service.erb
[Unit]
Description=Minecraft Server

Wants=network.target
After=network.target

[Service]
User=minecraft
WorkingDirectory=/home/minecraft
# This should look like ExecStart=/usr/bin/java -Xmx1009156 -jar server.jar
ExecStart=/usr/bin/java -Xmx<%= @memory_available %> -jar server.jar
ExecStop=/bin/kill -- $MAINPID
TimeoutStopSec=5
#
[Install]
WantedBy=multi-user.target
